# DEBIAN/prerm
#!/bin/sh
set -e

case "$1" in
    remove|upgrade)
        systemctl stop wg0.service || true
        systemctl disable wg0.service || true
    ;;
esac

exit 0
